# dockerfiles/framework/minio-nats/root/etc/s6/minio/run
#!/usr/bin/env bash

if [[ -f ./setup ]]; then
    ./setup
fi
exec minio server /export
