# dockerfiles/framework/minio-nats/root/etc/s6/nats/run
#!/usr/bin/env bash

if [[ -f ./setup ]]; then
    ./setup
fi
exec nats-streaming-server -m 8222
